# supermatrix/scripts/launchd/com.supermatrix.localwatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.supermatrix.localwatch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__TERMINAL_LAUNCHER__</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/terminal-launcher.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/terminal-launcher.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>
</dict>
</plist>
